# staged-systemd/forgejo.service.d/override.conf
[Unit]
After=network-online.target tailscaled.service
Wants=network-online.target tailscaled.service

[Service]
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do ip -4 addr show dev tailscale0 2>/dev/null | grep -q "100.94.10.36/32" && exit 0; sleep 1; done; echo "tailscale0 missing 100.94.10.36/32" >&2; exit 1'
