# staged-systemd/tower/vllm-backend.service
[Unit]
Description=DeepSeek-R1 backend (optional, port 8022)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/home/dino/bin/vllm-deepseek-r1-start.sh
Restart=on-failure
RestartSec=15
LimitMEMLOCK=infinity
LimitSTACK=67108864
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vllm-backend

[Install]
WantedBy=default.target
